# etcConfig/profile.d/spanish.sh
LANG=es_ES.UTF-8;						export LANG
CHARSET=UTF-8; 							export CHARSET

EDITOR='vim';							export EDITOR
PAGER='less';                                       		export PAGER

QT_QPA_PLATFORMTHEME='qt5ct';                       		export QT_QPA_PLATFORMTHEME
MOZ_ACCELERATED=1;						export MOZ_ACCELERATED
MOZ_WEBRENDER=1;						export MOZ_WEBRENDER

