# etcConfig/profile.d/spanish.sh
LANG=es_ES.UTF-8; export LANG
MM_CHARSET=UTF-8; export CHARSET
