# etcConfig/profile.d/spanish.sh
LANG="es_ES.UTF-8";						export LANG
CHARSET="UTF-8"; 						export CHARSET

EDITOR="vim";							export EDITOR
PAGER="less";                                       		export PAGER

QT_QPA_PLATFORMTHEME="qt5ct";                       		export QT_QPA_PLATFORMTHEME
MOZ_ACCELERATED=1;						export MOZ_ACCELERATED
MOZ_WEBRENDER=1;						export MOZ_WEBRENDER
XDG_RUNTIME_DIR="/var/run/user/$(id -u)";                       export XDG_RUNTIME_DIR
GOPATH="$HOME/.golang";		                                export GOPATH
RUST_SRC_PATH="/usr/local/lib/rustlib/src/rust/library";        export RUST_SRC_PATH
